# etc/systemd/system/qubes-update-check.service.d/40_qubes.conf
## This file is part of Qubes+Whonix.
## Copyright (C) 2015 Jason Mehring <nrgaway@gmail.com>
## License: GPL-2+

[Unit]
ConditionPathExists=!/var/run/qubes-service/whonix-template
